# vpower.service
[Unit]
Description=vpower daemon

[Service]
ExecStart=/usr/lib/vpower

[Install]
WantedBy=graphical.target
